# crontab
@reboot /home/ec2-user/setup.sh
*/5 * * * * /home/ec2-user/autoshutdown.sh